# view_updater/view_updater.cron
0 * * * * root /usr/local/bin/view_updater.py
